# deploy/stitch.service
# systemd unit for the Textile Stitch.
#
# This template is only for a root-managed Linux systemd service. Foreground
# operator runs use ~/Stitch/stitch.toml and ~/Stitch/stitch.env instead.
#
# Install:
#   sudo cp stitch /usr/local/bin/
#   sudo mkdir -p /etc/stitch-bot
#   sudo cp stitch.toml /etc/stitch-bot/stitch.toml
#   sudo cp stitch.key  /etc/stitch-bot/stitch.key   # holds the private key, chmod 600
#   sudo cp stitch.env  /etc/stitch-bot/stitch.env   # non-secret env, e.g. RUST_LOG
#   sudo chmod 600 /etc/stitch-bot/stitch.key
#   sudo cp stitch.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now stitch
#
# Logs:    journalctl -u stitch -f
# Restart: sudo systemctl restart stitch   (re-reads stitch.toml on start)

[Unit]
Description=Textile filler network operator bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# stitch.env is for non-secret settings only. The wallet key is loaded as a
# systemd credential and read from STITCH_PRIVATE_KEY_FILE.
EnvironmentFile=/etc/stitch-bot/stitch.env
LoadCredential=stitch.key:/etc/stitch-bot/stitch.key
Environment=STITCH_PRIVATE_KEY_FILE=%d/stitch.key
ExecStart=/usr/local/bin/stitch --config /etc/stitch-bot/stitch.toml
# The bot traps SIGTERM and exits cleanly after the current tick.
Restart=always
RestartSec=5
# Optional hardening.
DynamicUser=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadOnlyPaths=/etc/stitch-bot

[Install]
WantedBy=multi-user.target
